# init/sysdwitch.service
[Unit]
Description=SystemD Switch (sysdwitch)
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=/opt/sysdwitch
ExecStart=/opt/sysdwitch/sysdwitch
Restart=always
RestartSec=5

# Environment
EnvironmentFile=/opt/sysdwitch/configs/environments/local.env

# Security
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/sysdwitch

[Install]
WantedBy=default.target
