# configs/systemd/service-control.service
[Unit]
Description=Service Control Panel
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=%h/selfhost/apps/service-control-panel
ExecStart=%h/selfhost/apps/service-control-panel/service-control
Restart=always
RestartSec=5

# Environment
EnvironmentFile=%h/selfhost/apps/service-control-panel/configs/environments/local.env

# Security
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=%h/selfhost

[Install]
WantedBy=default.target
